# database/systemd/liqi-database-repository-check.timer
[Unit]
Description=Daily LIQI PostgreSQL repository verification

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
RandomizedDelaySec=10m
Persistent=true
Unit=liqi-database-repository-check.service

[Install]
WantedBy=timers.target
