# database/systemd/liqi-database-backup@.service
[Unit]
Description=LIQI PostgreSQL %i backup to independent TLS repository
Documentation=file:/opt/liqi/current/database/runbooks/backup-and-wal.md
After=postgresql-17.service liqi-secrets.service network-online.target
Wants=network-online.target
Requires=liqi-secrets.service
ConditionPathExists=/etc/liqi/database/backup.env

[Service]
Type=oneshot
User=postgres
Group=postgres
EnvironmentFile=/etc/liqi/database/backup.env
Environment=LIQI_SECRET_PGBACKREST_REPO_CA_FILE=%d/pgbackrest-repo-ca
Environment=LIQI_SECRET_PGBACKREST_REPO_CLIENT_CERT_FILE=%d/pgbackrest-repo-client-cert
Environment=LIQI_SECRET_PGBACKREST_REPO_CLIENT_KEY_FILE=%d/pgbackrest-repo-client-key
Environment=LIQI_SECRET_PGBACKREST_CIPHER_FILE=%d/pgbackrest-cipher-passphrase
ExecStart=/opt/liqi/current/database/bin/backup.sh %i
TimeoutStartSec=2h
UMask=0077
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUQuota=30%
MemoryMax=1536M
TasksMax=64
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
