# database/systemd/liqi-database-repository-check.service
[Unit]
Description=Check LIQI PostgreSQL archive and independent TLS repository
After=postgresql-17.service liqi-secrets.service network-online.target
Wants=network-online.target
Requires=liqi-secrets.service
ConditionPathExists=/etc/liqi/database/backup.env

[Service]
Type=oneshot
User=postgres
Group=postgres
EnvironmentFile=/etc/liqi/database/backup.env
Environment=LIQI_SECRET_PGBACKREST_REPO_CA_FILE=%d/pgbackrest-repo-ca
Environment=LIQI_SECRET_PGBACKREST_REPO_CLIENT_CERT_FILE=%d/pgbackrest-repo-client-cert
Environment=LIQI_SECRET_PGBACKREST_REPO_CLIENT_KEY_FILE=%d/pgbackrest-repo-client-key
Environment=LIQI_SECRET_PGBACKREST_CIPHER_FILE=%d/pgbackrest-cipher-passphrase
ExecStart=/opt/liqi/current/database/bin/repository-check.sh
TimeoutStartSec=30m
UMask=0077
CPUQuota=20%
MemoryMax=512M
TasksMax=32
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
