# infrastructure/systemd/liqi-database-credentials.service
[Unit]
Description=Restore approved LIQI database transient credentials
Requires=liqi-secrets.service postgresql-17.service
After=liqi-secrets.service postgresql-17.service
Before=pgbouncer.service liqi-beam.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/libexec/liqi-configure-database-credentials restore-transient --execute
TimeoutStartSec=90s
UMask=0077
RemainAfterExit=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER CAP_SETUID CAP_SETGID
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/run/liqi/secrets/database/database-role-passwords /var/lib/liqi/recovery/database-credentials-v1.json
ReadWritePaths=/run/liqi/secrets/database /run/liqi/secrets/beam

[Install]
WantedBy=multi-user.target
