# infrastructure/systemd/liqi-telemetry.slice
[Unit]
Description=LIQI bounded telemetry slice
Before=slices.target

[Slice]
Slice=liqi-platform.slice
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
CPUQuota=10%
CPUWeight=25
MemoryHigh=768M
MemoryMax=1G
MemorySwapMax=0
TasksMax=512
