# infrastructure/systemd/liqi-secrets.service
[Unit]
Description=Materialize LIQI service credentials from OCI Vault
After=network-online.target
Wants=network-online.target
Before=liqi-beam.service liqi-api.service liqi-realtime.service liqi-worker.service
ConditionPathExistsGlob=/etc/liqi/secrets/*.json

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/libexec/liqi-materialize-secrets --mapping-root /etc/liqi/secrets
TimeoutStartSec=120s
UMask=0077
RemainAfterExit=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/etc/liqi/secrets
ReadWritePaths=/run/liqi/secrets

[Install]
WantedBy=multi-user.target
