# infrastructure/systemd/liqi-edge.slice
[Unit]
Description=LIQI public edge slice
Before=slices.target

[Slice]
Slice=liqi-platform.slice
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
CPUQuota=10%
CPUWeight=50
MemoryHigh=192M
MemoryMax=256M
MemorySwapMax=0
TasksMax=512
