# infrastructure/systemd/liqi-beam.slice
[Unit]
Description=LIQI BEAM runtime slice
Before=slices.target

[Slice]
Slice=liqi-platform.slice
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
CPUQuota=180%
CPUWeight=200
MemoryHigh=7G
MemoryMax=8G
MemorySwapMax=0
TasksMax=4096
