# infrastructure/systemd/liqi-platform-runtime.slice
[Unit]
Description=LIQI disabled V0 compatibility runtime slice
Before=slices.target

[Slice]
Slice=liqi-platform.slice
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
CPUQuota=145%
CPUWeight=150
MemoryHigh=6G
MemoryMax=7G
MemorySwapMax=0
TasksMax=1536
