# infrastructure/systemd/liqi-host-readiness.service
[Unit]
Description=Validate LIQI V1 host runtime prerequisites
After=liqi-data-volume.service network-online.target liqi-secrets.service liqi-database-credentials.service postgresql-17.service pgbouncer.service caddy.service otelcol.service
Requires=liqi-data-volume.service liqi-secrets.service liqi-database-credentials.service postgresql-17.service pgbouncer.service caddy.service otelcol.service
Wants=network-online.target
Before=liqi-beam.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/liqi-host-readiness --output /run/liqi/host-runtime-v1.json
RemainAfterExit=yes
TimeoutStartSec=60s
UMask=0027
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/liqi

[Install]
WantedBy=multi-user.target
